# etc/systemd.service
[Unit]
Description=__NAME__ Server
DefaultDependencies=true

[Service]
EnvironmentFile="__PDIR__/.env"
ExecStart="__PDIR__/.local/bin/gunicorn project.wsgi:application"
Restart=on−failure

[Install]
WantedBy=multi-user.target
